# deployment/ansible/ddos-mitigation.service.j2
[Unit]
Description=DDoS Protection System - Mitigation
After=network.target
Wants=network.target

[Service]
Type=simple
User={{ service_user }}
Group={{ service_group }}
WorkingDirectory={{ install_dir }}
Environment="PATH={{ install_dir }}/venv/bin"
Environment="PYTHONPATH={{ install_dir }}"
ExecStart={{ install_dir }}/venv/bin/python {{ install_dir }}/scripts/run_mitigation.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
